# DEBIAN/postinst
#!/bin/bash

# Reload systemd dan aktifkan unit terkait
systemctl daemon-reexec
systemctl daemon-reload
systemctl enable knot-init.service
systemctl enable knot-delay.service
systemctl enable kresd-rpz-update.timer

# Pastikan /dev/shm diatur ke 24G
if ! grep -q "^tmpfs[[:space:]]\+/dev/shm" /etc/fstab; then
    echo "tmpfs /dev/shm tmpfs defaults,size=24G 0 0" >> /etc/fstab
else
    sed -i 's|^tmpfs[[:space:]]\+/dev/shm.*|tmpfs /dev/shm tmpfs defaults,size=24G 0 0|' /etc/fstab
fi

# Remount untuk menerapkan ukuran baru
mount -o remount /dev/shm || true

# Buat direktori RAM dan beri izin
mkdir -p /dev/shm/kresd-cache
mkdir -p /dev/shm/kresd-rpz
chown knot-resolver:knot-resolver /dev/shm/kresd-cache
chown knot-resolver:knot-resolver /dev/shm/kresd-rpz

# Buat symlink RPZ agar terbaca oleh Knot Resolver
ln -sf /dev/shm/kresd-rpz/xdomains.rpz   /etc/knot-resolver/xdomains.rpz
ln -sf /dev/shm/kresd-rpz/safesearch.rpz /etc/knot-resolver/safesearch.rpz
ln -sf /dev/shm/kresd-rpz/custom.rpz     /etc/knot-resolver/custom.rpz
